# extractor/crontab
0 5 * * * /usr/bin/docker stop $(docker ps -a -q); /usr/bin/docker run -d -p 5000:5000 extractor/flask-app
